# config/crontab
SHELL=/bin/bash
* * * * * if [[ -f ~/bin/keys.bash ]]; then ~/bin/keys.bash > /tmp/keys.log; fi
