# config/systemd/authorized_keys.timer
[Timer]
OnActiveSec=60
RandomizedDelaySec=15
OnCalendar=*-*-* *:30,00:00
Unit=authorized_keys.service

[Install]
WantedBy=default.target
